# cron/crontab
*/1 *   *   *   *   XDG_RUNTIME_DIR=/run/user/$(id -u) /home/albert/bin/battery_level
